# collect_kernel_logs.timer
[Unit]
Description=Collect Kernel Logs Timer

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target

